# dist/cloud-collector/cloud-collector.service
[Unit]
Description=OpenTelemetry-eBPF Cloud Collector
After=network.target

[Service]
Type=exec
ExecStart=/bin/bash -c 'exec /usr/bin/cloud-collector $(< /etc/opentelemetry-ebpf/cloud-collector.args)'
Restart=on-failure

[Install]
WantedBy=multi-user.target
